# contrib/init/fartcoind.conf
description "Fartcoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env FARTCOIND_BIN="/usr/bin/fartcoind"
env FARTCOIND_USER="fartcoin"
env FARTCOIND_GROUP="fartcoin"
env FARTCOIND_PIDDIR="/var/run/fartcoind"
# upstart can't handle variables constructed with other variables
env FARTCOIND_PIDFILE="/var/run/fartcoind/fartcoind.pid"
env FARTCOIND_CONFIGFILE="/etc/fartcoin/fartcoin.conf"
env FARTCOIND_DATADIR="/var/lib/fartcoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # fartcoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$FARTCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run fartcoind."
        echo "The setting must appear in $FARTCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Fartcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$FARTCOIND_PIDDIR"
    chmod 0755 "$FARTCOIND_PIDDIR"
    chown $FARTCOIND_USER:$FARTCOIND_GROUP "$FARTCOIND_PIDDIR"
    chown $FARTCOIND_USER:$FARTCOIND_GROUP "$FARTCOIND_CONFIGFILE"
    chmod 0660 "$FARTCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$FARTCOIND_PIDFILE" \
    --chuid $FARTCOIND_USER:$FARTCOIND_GROUP \
    --exec "$FARTCOIND_BIN" \
    -- \
    -pid="$FARTCOIND_PIDFILE" \
    -conf="$FARTCOIND_CONFIGFILE" \
    -datadir="$FARTCOIND_DATADIR" \
    -disablewallet \
    -daemon

